# camoufox/supervisor/conf.d/fluxbox.conf
[program:fluxbox]
command=fluxbox
autorestart=true
stdout_logfile=/app/supervisor/logs/fluxbox.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
redirect_stderr=true
